# config/deploy/templates/shared/log_forwarder.service.erb
[Unit]
Description=<%= fetch(:application) %>
After=network-online.target

[Service]
User=resfinity
EnvironmentFile=/etc/environment
Environment=GOOGLE_APPLICATION_CREDENTIALS=/opt/wheels/current/anixe-extra.json
Environment=ASPNETCORE_ENVIRONMENT=<%= fetch(:stage) %>
WorkingDirectory=<%= fetch(:deploy_to) %>/current/
ExecStart=/usr/share/dotnet/dotnet <%= fetch(:application) %>.dll -p /var/log/resfinity/wheels/ -f .complete -s scripts/wheels_gcs.cs.txt -w 4 -b gcs -v
Restart=always

[Install]
WantedBy=multi-user.target